# squashfs-root-0/usr/lib/opkg/info/cwmp.list
/usr/bin/tr069/lan_ipv6
/usr/bin/tr069/myddns
/usr/bin/tr069/myportforward
/etc/hotplug.d/firewall/98-cwmp
/usr/bin/tr069/myrestore
/usr/bin/tr069/mybackup
/usr/bin/cwmp
/etc/init.d/cwmp
/usr/bin/tr069/wan_ipv6
/usr/bin/tr069/check_device_sn
/usr/bin/tr069/myremote
/etc/tr098_ap.xml
/usr/bin/tr069/myupgrade
/usr/bin/tr069/myppp
/etc/hotplug.d/iface/98-cwmp_hotplug
/usr/bin/tr069/mypasswd
/usr/lib/lua/luci/controller/admin/cwmp.lua
/usr/lib/libtr098.so
/etc/tr098.xml
/usr/bin/tr069/diagnostic
/usr/bin/cwmp_cli
/usr/bin/tr069/generate_hosts
/usr/bin/tr069/get_wifimac
/usr/bin/tr069/check_params
/usr/bin/tr069/mynat
