# squashfs-root-0/usr/lib/opkg/info/kmod-shortcut-fe.list
/etc/modules.d/47-shortcut-fe
/lib/modules/iplatform/shortcut-fe.ko
/lib/modules/iplatform/shortcut-fe-ipv6.ko
/etc/init.d/shortcut-fe
