# squashfs-root-0/usr/lib/opkg/info/wireless-script-mtk.list
/lib/wifi/tplink_mtk.sh
/sbin/is_cal.lua
/lib/wifi/triband/wlan.ini
/lib/wifi/dualband/wlan.ini
/etc/init.d/wireless
/sbin/wifi_button
/lib/wifi/wps_button
/sbin/is_btn_test.lua
/sbin/wifi_firm
/etc/rc.d/S30wireless
/lib/wifi/mtk_wlan_var.sh
/lib/wifi/config.ini
/lib/wifi/mtk_wlan_extend.sh
/sbin/is_btn_test
/sbin/is_cal
