# squashfs-root/usr/lib/opkg/info/autodial.list
/etc/hotplug.d/cable/00-autodial
/etc/rc.d/K10reboot
/lib/autodial/auto_dial.lua
