# squashfs-root/usr/lib/opkg/info/network-roaming.list
/etc/init.d/nrd
/lib/functions/nrd-config.sh
/usr/sbin/nrd
/etc/hotplug.d/mode/20-nrd
