# squashfs-root/usr/lib/opkg/info/product-file.list
/lib/preinit/06_init_insmod_usbhost_scsi
/usr/sbin/wps_status
/usr/sbin/status
/etc/init.d/hnat-start
/lib/led_switch
/lib/preinit/71_init_etc_ramfs
/usr/sbin/lan_connected
/usr/sbin/wan_connected
/usr/sbin/sync_hnat_status
