# squashfs-root/usr/lib/opkg/info/yandex-dns.list
/lib/upgrade/keep.d/yandex-dns
/lib/yandex_dns/core_log.sh
/etc/init.d/yandex_dns
/etc/config/yandex_dns
/lib/yandex_dns/core.sh
/lib/yandex_dns/core_yandex_dns.sh
/etc/hotplug.d/firewall/22-yandex_dns
/lib/yandex_dns/config.sh
